# cluster-node/rootfs/etc/s6-overlay/s6-rc.d/corosync/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Cluster Node
# Runs the Corosync daemon
# ==============================================================================

# # Wait for PHP-FPM to become available
bashio::net.wait_for 5403
bashio::net.wait_for 5403
bashio::net.wait_for 5403

bashio::log.info "Starting Corosync...."

exec corosync -f